# scripts/com.stoneforge.weekly-retro.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stoneforge.weekly-retro</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/bun</string>
        <string>run</string>
        <string>scripts/retro.ts</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/rahi-akila/stoneforge</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/stoneforge-retro.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/stoneforge-retro.err</string>
</dict>
</plist>
